# root/etc/services.d/paper/run
#!/usr/bin/with-contenv bash

cd /data || exit

exec \
	s6-setuidgid abc java \
    -Xms2G \
    -Xmx2G \
	-jar /data/paper.jar nogui
