# root/etc/services.d/paper/run
#!/usr/bin/with-contenv bash

cd /data || exit

exec \
	s6-setuidgid abc java \
    -Xms6G \
    -Xmx6G \
	-jar /data/paper.jar nogui
